# backend/scripts/deploy/templates/systemd/eyefi-nest-api-b.service
[Unit]
Description=Eyefi Nest API Service (instance 2)
After=network.target

[Service]
Type=simple
User=rda
WorkingDirectory=/var/www/html/backend/nest-api
EnvironmentFile=/var/www/html/backend/.env.production
Environment=PORT=3003
Environment=NODE_ENV=production
Environment=ODBCINI=/var/www/html/backend/qad-driver/Connect64_for_ODBC_71/odbc.ini
Environment=ODBCINST=/var/www/html/backend/qad-driver/Connect64_for_ODBC_71/odbcinst.ini
Environment=LD_LIBRARY_PATH=/var/www/html/backend/qad-driver/Connect64_for_ODBC_71/lib
ExecStart=/usr/bin/env PORT=3003 /home/rda/.nvm/versions/node/v20.20.2/bin/node /var/www/html/backend/nest-api/dist/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
